# roles/docker_prune/templates/docker_prune.service.j2
[Unit]
Description=Prune orphaned docker primitives
[Service]
Type=simple
ExecStart=/usr/bin/docker {{ docker_prune_type }} prune --all --force
